# backend/reverb.service
[Unit]
Description=Reverb service to control websocket server for agile-tools.dk
Conflicts=shutdown.target

[Service]
User=www-data
ExecStart=/usr/bin/php artisan reverb:start
WorkingDirectory=/srv/sites/agile-tools.dk/backend
Restart=on-failure
RestartSec=20

[Install]
WantedBy=multi-user.target
